# supervisor-bitcoin.conf
[supervisord]
nodaemon=true

[program:bitcoind]
command=su -c "bitcoind -conf=/home/crypto/bitcoin.conf" crypto
autorestart=true
